# debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ] ; then
    # Reload the systemd user daemon for all users who are currently
    # logged in, to make them aware of the new units.
    # This is a best-effort approach. A user logging in later will
    # see the units automatically.
    for UID_PATH in /run/user/*; do
        if [ -d "$UID_PATH" ] && [ -S "$UID_PATH/bus" ]; then
            UID=$(basename "$UID_PATH")
            # Ensure UID is a number to avoid issues with non-user directories
            if [ -n "$UID" ] && [ "$UID" -eq "$UID" ] 2>/dev/null; then
                echo "Processing user $UID..."
                HOME_DIR=$(getent passwd "$UID" | cut -d: -f6)
                if [ -n "$HOME_DIR" ] && [ -d "$HOME_DIR" ]; then
                    export XDG_RUNTIME_DIR="/run/user/$UID"
                    export DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$UID/bus"

                    # Check if systemd user instance is running
                    if sudo -u "#$UID" systemctl --user is-system-running >/dev/null 2>&1; then
                        echo "Reloading systemd --user daemon for UID $UID"
                        sudo -u "#$UID" systemctl --user daemon-reload || echo "Failed to reload daemon for $UID, continuing..."

                        echo "Enabling fluxfce services for UID $UID"
                        sudo -u "#$UID" systemctl --user enable fluxfce-scheduler.timer fluxfce-login.service fluxfce-resume.service || echo "Failed to enable services for $UID, continuing..."
                    else
                        echo "Systemd user instance not running for UID $UID or not accessible."
                    fi
                else
                    echo "Could not determine home directory for UID $UID or directory does not exist."
                fi
            fi
        fi
    done
fi

exit 0
